# package/lean/luci-app-wiremesh/root/etc/uci-defaults/luci-wiremesh
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@wiremesh[-1]
	add ucitrack wiremesh
	set ucitrack.@wiremesh[-1].init=wiremesh
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
